# deploy/smart_edu.service
[Unit]
Description=Smart Edu FastAPI Service
After=network.target mysql.service neo4j.service

[Service]
Type=simple
User=smart_edu
Group=smart_edu
WorkingDirectory=/opt/smart_edu/src
EnvironmentFile=/opt/smart_edu/.env
ExecStart=/opt/smart_edu/venv/bin/python -m backend.app
Restart=on-failure
RestartSec=5s

# 日志输出到 systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smart_edu

[Install]
WantedBy=multi-user.target
